# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:polymarket-bot]
command=python -m src.main
directory=/app
autostart=true
autorestart=true
startretries=10
stderr_logfile=/app/logs/bot_error.log
stdout_logfile=/app/logs/bot_output.log
environment=PYTHONUNBUFFERED=1
